# install-avahi-commands
# Disable Avahi
sudo systemctl stop avahi-daemon.socket avahi-daemon.service avahi-daemon
sudo systemctl disable avahi-daemon.socket avahi-daemon.service avahi-daemon
sudo apt purge -y avahi-daemon

# If Purge Not Purging Avahi

### Locate All Files and Directories Containing "avahi"
sudo find / -path "/home" -prune -o -name "*avahi*" -print 2>/dev/null 

### Delete All Found Files and Directories
sudo find / -path "/home/pi" -prune -o -name "*avahi*" -exec rm -rf {} + 2>/dev/null 

### Remove Any avahi Packages
sudo apt purge avahi* -y


# Switch apt-get method
sudo cp /etc/apt/sources.list /etc/apt/sources.list~

# Remove "#" infront of deb-src
sudo nano /etc/apt/sources.list

# Do if removing "#" doesn't work
sudo sed -Ei 's/#deb-src /deb-src /' /etc/apt/sources.list

# For Ubuntu:
# sed -Ei 's/# deb-src /deb-src /' /etc/apt/sources.list
sudo apt-get update

# Install Dependencies
sudo apt-get build-dep -y avahi
sudo apt install -y libevent-dev qtbase5-dev mono-mcs monodoc-http

# Build Avahi
sudo ./bootstrap.sh
sudo make
sudo make install

# Find Libraries
sudo /sbin/ldconfig -v

# Add avahi user
sudo adduser --disabled-password --quiet --system --home /var/run/avahi-daemon --gecos "Avahi mDNS daemon" --group avahi

# Start avahi
sudo systemctl daemon-reload
sudo systemctl start avahi-daemon
sudo systemctl enable avahi-daemon

***

# If Systemctl Enable Not Working
sudo nano /etc/init.d/avahi-daemon

## add: 
### BEGIN INIT INFO
# Default-Start: 2 3 4 5
### END INIT INFO

***

# Reboot Time
sudo reboot

# Check If Changed
netstat -lntu | grep 5353 ### Should have no output
netstat -lntu | grep 5358 ### Should have a output
